# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/13/bin/postgres -D /var/lib/postgresql/13/main -c config_file=/etc/postgresql/13/main/postgresql.conf
user=postgres
autorestart=true
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.err

[program:redis]
command=redis-server --daemonize no
user=redis
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.err

[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address ":9001" --address ":9000"
user=root
autorestart=true
stdout_logfile=/var/log/minio.log
stderr_logfile=/var/log/minio.err

[program:backend]
directory=/app/backend
command=/opt/venv/backend/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
environment=PYTHONPATH="/app/backend"
user=root
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.err

[program:celery_worker]
directory=/app/backend
command=/opt/venv/backend/bin/celery -A app.core.celery_app worker --loglevel=info
environment=PYTHONPATH="/app/backend"
user=root
autorestart=true
stdout_logfile=/var/log/celery.log
stderr_logfile=/var/log/celery.err

[program:chatbot]
directory=/app/chatbot
command=/opt/venv/chatbot/bin/rasa run -m models --enable-api --cors "*" --port 5005 --debug
user=root
autorestart=true
stdout_logfile=/var/log/chatbot.log
stderr_logfile=/var/log/chatbot.err

[program:chatbot_actions]
directory=/app/chatbot
command=/opt/venv/chatbot/bin/rasa run actions --port 5055
user=root
autorestart=true
stdout_logfile=/var/log/chatbot_actions.log
stderr_logfile=/var/log/chatbot_actions.err

[program:nginx]
command=nginx -g 'daemon off;'
user=root
autorestart=true
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
